# roles/restic/templates/restic-backup.service.j2
[Unit]
Description=Restic backup on %I
After=network-online.target

[Service]
User=restic
Group=restic
Type=oneshot
EnvironmentFile=/etc/restic/%I.env
