# roles/restic/templates/restic-backup.dropin.timer.j2
[Timer]
OnCalendar={{ backup_config.timer | default(backup_timer) }}
