# roles/restic/templates/restic-backup.dropin.service.j2
[Service]
{% for command in backup_commands -%}
ExecStart={{ command }}
{% endfor %}
AmbientCapabilities=CAP_DAC_READ_SEARCH
