# roles/restic/templates/restic-forget.dropin.service.j2
[Service]
{% for command in backup_forget_commands -%}
ExecStart={{ command }}
{% endfor %}
