# roles/restic/templates/restic-forget.dropin.timer.j2
[Timer]
OnCalendar={{ backup_config.forget_timer | default('') }}
